# deploy/systemd/kcml.service
[Unit]
Description=KCML MCP Server Manager
After=network.target postgresql.service

[Service]
Type=simple
User=kcml
Group=kcml
WorkingDirectory=/opt/kcml/current
EnvironmentFile=/etc/kcml/kcml.env
ExecStart=/usr/bin/node apps/server/dist/index.js
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/kcml /var/log/kcml /var/lib/kcml/runtime

[Install]
WantedBy=multi-user.target
