# deploy/systemd/kcml-onboarding-worker.service
[Unit]
Description=KCML automated onboarding worker
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=kcml
Group=kcml
WorkingDirectory=/opt/kcml/current
EnvironmentFile=/etc/kcml/worker.env
Environment=KCML_PROCESS_ROLE=worker
Environment=DATABASE_URL_FILE=%d/database_url
Environment=EGRESS_CAPABILITY_HMAC_KEY_BASE64_FILE=%d/egress_capability_hmac
Environment=GITHUB_TOKEN_FILE=%d/github_token
Environment=GITHUB_APP_PRIVATE_KEY_BASE64_FILE=%d/github_app_private_key
LoadCredential=database_url:/etc/kcml/credentials/worker/database_url
LoadCredential=egress_capability_hmac:/etc/kcml/credentials/worker/egress_capability_hmac
LoadCredential=github_token:/etc/kcml/credentials/worker/github_token
LoadCredential=github_app_private_key:/etc/kcml/credentials/worker/github_app_private_key
Environment=HOME=/var/lib/kcml/podman
Environment=XDG_DATA_HOME=/var/lib/kcml/podman/data
Environment=XDG_CONFIG_HOME=/var/lib/kcml/podman/config
# A trusted production runner may stage short-lived GHCR credentials here.
# Both files live below the protected Podman home and are removed as soon as
# the worker has verified the signed artifact.
Environment=REGISTRY_AUTH_FILE=/var/lib/kcml/podman/auth.json
Environment=DOCKER_CONFIG=/var/lib/kcml/podman/.docker
ExecStart=/usr/bin/node apps/server/dist/cli/onboarding-worker.js
Restart=on-failure
RestartSec=5
# Podman is submitted to the persistent kcml user manager, so the hardened
# worker never executes newuidmap/newgidmap itself. Handler containers retain
# their own cap-drop and no-new-privileges restrictions.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
ReadWritePaths=/var/lib/kcml/onboarding /var/lib/kcml/runtime /var/lib/kcml/podman

[Install]
WantedBy=multi-user.target
