# deploy/systemd/kcml-monitor.service
[Unit]
Description=KCML recertification and MCP monitor
After=network-online.target postgresql.service kcml.service
Wants=network-online.target

[Service]
Type=simple
User=kcml
Group=kcml
WorkingDirectory=/opt/kcml/current
EnvironmentFile=/etc/kcml/monitor.env
Environment=KCML_PROCESS_ROLE=monitor
Environment=DATABASE_URL_FILE=%d/database_url
Environment=EGRESS_CAPABILITY_HMAC_KEY_BASE64_FILE=%d/egress_capability_hmac
Environment=ALERT_PRIMARY_HMAC_KEY_BASE64_FILE=%d/alert_primary_hmac
Environment=ALERT_BACKUP_HMAC_KEY_BASE64_FILE=%d/alert_backup_hmac
LoadCredential=database_url:/etc/kcml/credentials/monitor/database_url
LoadCredential=egress_capability_hmac:/etc/kcml/credentials/monitor/egress_capability_hmac
LoadCredential=alert_primary_hmac:/etc/kcml/credentials/monitor/alert_primary_hmac
LoadCredential=alert_backup_hmac:/etc/kcml/credentials/monitor/alert_backup_hmac
Environment=HOME=/var/lib/kcml/podman
Environment=XDG_DATA_HOME=/var/lib/kcml/podman/data
Environment=XDG_CONFIG_HOME=/var/lib/kcml/podman/config
Environment=XDG_RUNTIME_DIR=/run/kcml-monitor-podman
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/kcml-monitor-podman/bus
RuntimeDirectory=kcml-monitor-podman
RuntimeDirectoryMode=0700
ExecStart=/usr/bin/node apps/server/dist/cli/monitor.js
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
ReadWritePaths=/var/lib/kcml/runtime /var/lib/kcml/podman

[Install]
WantedBy=multi-user.target
